# deploy/cloudflared-home.service
[Unit]
Description=Cloudflare Tunnel — home-lettuce.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vantop
ExecStart=/home/vantop/.local/bin/cloudflared tunnel --config /home/vantop/.cloudflared/home-lettuce.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
