# deploy/vantop-dashboard.service
[Unit]
Description=Vantop Dashboard (Flask)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vantop
WorkingDirectory=/home/vantop/HomeLettuce/Vantop_Dashboard
ExecStart=/home/vantop/HomeLettuce/Vantop_Dashboard/venv/bin/python3 app.py
Restart=always
RestartSec=5

# El .env se lee desde el WorkingDirectory via python-dotenv
# Si necesitas variables extra, agrégalas aquí:
# Environment=VARIABLE=valor

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/vantop/vantop-token.json /home/vantop/HomeLettuce/Vantop_Dashboard
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
